# runner.service
[Unit]
Description=myService 

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/